# .config/prod/supervisord.conf
[program:nginx]
command = nginx -g 'daemon off;'

[program:uwsgi]
command = uwsgi -i /srv/project/.config/prod/uwsgi.ini

[program:vuejs]
command = npm run dev --prefix /srv/frontend